# os/etc/cron.hourly/logrotate.sh
#!/bin/bash

logrotate --skip-state-lock --force /etc/logrotate.d/app
